# DEBIAN/postinst
#!/bin/bash

systemctl enable cups-printer-http
systemctl start cups-printer-http
